# ops/systemd/nba-retrain.timer
[Unit]
Description=Run NBA V2 Retraining Weekly

[Timer]
OnCalendar=Mon *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
