# ops/systemd/nhl-totals-kpi.service
[Unit]
Description=NHL Totals V2 - KPI Report
After=docker.service network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/Philly-P-Sniper
ExecStart=/usr/bin/env bash -lc 'cd /home/ubuntu/Philly-P-Sniper && docker exec -i -e NHL_TOTALS_V2_ENABLED=true philly_p_api python3 scripts/ops/generate_nhl_kpi.py'
StandardOutput=append:/home/ubuntu/Philly-P-Sniper/logs/systemd/nhl-totals-kpi.log
StandardError=inherit

[Install]
WantedBy=multi-user.target
