# ops/systemd/nhl-odds-ingest.timer
[Unit]
Description=Timer for NHL Odds Ingest (Every 60m)

[Timer]
OnCalendar=*:00:00
Unit=nhl-odds-ingest.service

[Install]
WantedBy=timers.target
